# crontab
# cron job for the FreePWING conversions
#
# m h dom mon dow  command
03 18   *   *   5  $HOME/usr/src/fpw-scripts/convert-all.sh > $HOME/fpw-convert.log 2>&1
